# debian/python-seed-cleaner.cron.d

0   */2     *   *   *   root    /usr/bin/seedcleaner.py

