# debian/python-django-seed.prerm
#!/bin/bash

if [[ -x "/etc/init.d/uwsgi" ]]; then
    invoke-rc.d uwsgi stop
fi

if [[ -x "/etc/init.d/nginx" ]]; then
    invoke-rc.d nginx stop
fi

